# ports/multimedia/subsonic/pkg-plist
@unexec %D/etc/rc.d/%%PORTNAME%% forcestop 2>/dev/null || true
%%TRANSCODING%%%%ETCDIR%%/transcode/ffmpeg
%%TRANSCODING%%@dirrmtry %%ETCDIR%%/transcode
%%DATADIR%%/Subsonic/subsonic-booter-jar-with-dependencies.jar
%%DATADIR%%/Subsonic/subsonic.war
@dirrmtry %%DATADIR%%/Subsonic
@dirrmtry %%DATADIR%%
@unexec rmdir /var/db/%%PORTNAME%%
@exec mkdir -p /var/db/%%PORTNAME%%
@exec chown -R %%USER%%:%%GROUP%% /var/db/%%PORTNAME%%
@sample etc/%%PORTNAME%%.conf.sample
