# ports/multimedia/sonarr/files/sonarr.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: %%PORTNAME%%
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# %%PORTNAME%%_enable:	Set to YES to enable %%PORTNAME%%
#			Default: NO
# %%PORTNAME%%_user:	The user account used to run the %%PORTNAME%% daemon.
#			This is optional, however do not specifically set this to an
#			empty string as this will cause the daemon to run as root.
#			Default: %%USER%%
# %%PORTNAME%%_group:	The group account used to run the %%PORTNAME%% daemon.
#			This is optional, however do not specifically set this to an
#			empty string as this will cause the daemon to run with group wheel.
#			Default: %%GROUP%%
# %%PORTNAME%%_data_dir:	Directory where %%PORTNAME%% configuration
#			data is stored.
#			Default: /var/db/%%PORTNAME%%

. /etc/rc.subr
name=%%PORTNAME%%
rcvar=${name}_enable
load_rc_config $name

: ${%%PORTNAME%%_enable:="NO"}
: ${%%PORTNAME%%_user:="%%USER%%"}
: ${%%PORTNAME%%_group:="%%GROUP%%"}
: ${%%PORTNAME%%_data_dir:="/var/db/%%PORTNAME%%"}

pidfile="${%%PORTNAME%%_data_dir}/nzbdrone.pid"
command="/usr/sbin/daemon"
procname="%%LOCALBASE%%/bin/mono"
command_args="-f ${procname} %%DATADIR%%/NzbDrone/NzbDrone.exe --data=${%%PORTNAME%%_data_dir} --nobrowser"

start_precmd=%%PORTNAME%%_precmd
%%PORTNAME%%_precmd() {
	if [ ! -d ${%%PORTNAME%%_data_dir} ]; then
		install -d -o ${%%PORTNAME%%_user} -g ${%%PORTNAME%%_group} ${%%PORTNAME%%_data_dir}
	fi
}

run_rc_command "$1"
