# ports/games/gamez/pkg-plist
@unexec %D/etc/rc.d/gamez forcestop 2>/dev/null || true
%%DATADIR%%/Gamez/API.md
%%DATADIR%%/Gamez/Gamez.py
%%DATADIR%%/Gamez/Gamez.pyproj
%%DATADIR%%/Gamez/Gamez.sln
%%DATADIR%%/Gamez/Gamez.suo
%%DATADIR%%/Gamez/README.md
%%DATADIR%%/Gamez/cherrypy/LICENSE.txt
%%DATADIR%%/Gamez/cherrypy/__init__.py
%%DATADIR%%/Gamez/cherrypy/_cpchecker.py
%%DATADIR%%/Gamez/cherrypy/_cpcompat.py
%%DATADIR%%/Gamez/cherrypy/_cpconfig.py
%%DATADIR%%/Gamez/cherrypy/_cpdispatch.py
%%DATADIR%%/Gamez/cherrypy/_cperror.py
%%DATADIR%%/Gamez/cherrypy/_cplogging.py
%%DATADIR%%/Gamez/cherrypy/_cpmodpy.py
%%DATADIR%%/Gamez/cherrypy/_cpnative_server.py
%%DATADIR%%/Gamez/cherrypy/_cpreqbody.py
%%DATADIR%%/Gamez/cherrypy/_cprequest.py
%%DATADIR%%/Gamez/cherrypy/_cpserver.py
%%DATADIR%%/Gamez/cherrypy/_cpthreadinglocal.py
%%DATADIR%%/Gamez/cherrypy/_cptools.py
%%DATADIR%%/Gamez/cherrypy/_cptree.py
%%DATADIR%%/Gamez/cherrypy/_cpwsgi.py
%%DATADIR%%/Gamez/cherrypy/_cpwsgi_server.py
%%DATADIR%%/Gamez/cherrypy/cherryd
%%DATADIR%%/Gamez/cherrypy/favicon.ico
%%DATADIR%%/Gamez/cherrypy/lib/__init__.py
%%DATADIR%%/Gamez/cherrypy/lib/auth.py
%%DATADIR%%/Gamez/cherrypy/lib/auth_basic.py
%%DATADIR%%/Gamez/cherrypy/lib/auth_digest.py
%%DATADIR%%/Gamez/cherrypy/lib/caching.py
%%DATADIR%%/Gamez/cherrypy/lib/covercp.py
%%DATADIR%%/Gamez/cherrypy/lib/cpstats.py
%%DATADIR%%/Gamez/cherrypy/lib/cptools.py
%%DATADIR%%/Gamez/cherrypy/lib/encoding.py
%%DATADIR%%/Gamez/cherrypy/lib/gctools.py
%%DATADIR%%/Gamez/cherrypy/lib/http.py
%%DATADIR%%/Gamez/cherrypy/lib/httpauth.py
%%DATADIR%%/Gamez/cherrypy/lib/httputil.py
%%DATADIR%%/Gamez/cherrypy/lib/jsontools.py
%%DATADIR%%/Gamez/cherrypy/lib/profiler.py
%%DATADIR%%/Gamez/cherrypy/lib/reprconf.py
%%DATADIR%%/Gamez/cherrypy/lib/sessions.py
%%DATADIR%%/Gamez/cherrypy/lib/static.py
%%DATADIR%%/Gamez/cherrypy/lib/xmlrpcutil.py
%%DATADIR%%/Gamez/cherrypy/process/__init__.py
%%DATADIR%%/Gamez/cherrypy/process/plugins.py
%%DATADIR%%/Gamez/cherrypy/process/servers.py
%%DATADIR%%/Gamez/cherrypy/process/win32.py
%%DATADIR%%/Gamez/cherrypy/process/wspbus.py
%%DATADIR%%/Gamez/cherrypy/scaffold/__init__.py
%%DATADIR%%/Gamez/cherrypy/scaffold/apache-fcgi.conf
%%DATADIR%%/Gamez/cherrypy/scaffold/example.conf
%%DATADIR%%/Gamez/cherrypy/scaffold/site.conf
%%DATADIR%%/Gamez/cherrypy/scaffold/static/made_with_cherrypy_small.png
%%DATADIR%%/Gamez/cherrypy/test/__init__.py
%%DATADIR%%/Gamez/cherrypy/test/_test_decorators.py
%%DATADIR%%/Gamez/cherrypy/test/_test_states_demo.py
%%DATADIR%%/Gamez/cherrypy/test/benchmark.py
%%DATADIR%%/Gamez/cherrypy/test/checkerdemo.py
%%DATADIR%%/Gamez/cherrypy/test/helper.py
%%DATADIR%%/Gamez/cherrypy/test/logtest.py
%%DATADIR%%/Gamez/cherrypy/test/modfastcgi.py
%%DATADIR%%/Gamez/cherrypy/test/modfcgid.py
%%DATADIR%%/Gamez/cherrypy/test/modpy.py
%%DATADIR%%/Gamez/cherrypy/test/modwsgi.py
%%DATADIR%%/Gamez/cherrypy/test/sessiondemo.py
%%DATADIR%%/Gamez/cherrypy/test/static/dirback.jpg
%%DATADIR%%/Gamez/cherrypy/test/static/index.html
%%DATADIR%%/Gamez/cherrypy/test/style.css
%%DATADIR%%/Gamez/cherrypy/test/test.pem
%%DATADIR%%/Gamez/cherrypy/test/test_auth_basic.py
%%DATADIR%%/Gamez/cherrypy/test/test_auth_digest.py
%%DATADIR%%/Gamez/cherrypy/test/test_bus.py
%%DATADIR%%/Gamez/cherrypy/test/test_caching.py
%%DATADIR%%/Gamez/cherrypy/test/test_config.py
%%DATADIR%%/Gamez/cherrypy/test/test_config_server.py
%%DATADIR%%/Gamez/cherrypy/test/test_conn.py
%%DATADIR%%/Gamez/cherrypy/test/test_core.py
%%DATADIR%%/Gamez/cherrypy/test/test_dynamicobjectmapping.py
%%DATADIR%%/Gamez/cherrypy/test/test_encoding.py
%%DATADIR%%/Gamez/cherrypy/test/test_etags.py
%%DATADIR%%/Gamez/cherrypy/test/test_http.py
%%DATADIR%%/Gamez/cherrypy/test/test_httpauth.py
%%DATADIR%%/Gamez/cherrypy/test/test_httplib.py
%%DATADIR%%/Gamez/cherrypy/test/test_json.py
%%DATADIR%%/Gamez/cherrypy/test/test_logging.py
%%DATADIR%%/Gamez/cherrypy/test/test_mime.py
%%DATADIR%%/Gamez/cherrypy/test/test_misc_tools.py
%%DATADIR%%/Gamez/cherrypy/test/test_objectmapping.py
%%DATADIR%%/Gamez/cherrypy/test/test_proxy.py
%%DATADIR%%/Gamez/cherrypy/test/test_refleaks.py
%%DATADIR%%/Gamez/cherrypy/test/test_request_obj.py
%%DATADIR%%/Gamez/cherrypy/test/test_routes.py
%%DATADIR%%/Gamez/cherrypy/test/test_session.py
%%DATADIR%%/Gamez/cherrypy/test/test_sessionauthenticate.py
%%DATADIR%%/Gamez/cherrypy/test/test_states.py
%%DATADIR%%/Gamez/cherrypy/test/test_static.py
%%DATADIR%%/Gamez/cherrypy/test/test_tools.py
%%DATADIR%%/Gamez/cherrypy/test/test_tutorials.py
%%DATADIR%%/Gamez/cherrypy/test/test_virtualhost.py
%%DATADIR%%/Gamez/cherrypy/test/test_wsgi_ns.py
%%DATADIR%%/Gamez/cherrypy/test/test_wsgi_vhost.py
%%DATADIR%%/Gamez/cherrypy/test/test_wsgiapps.py
%%DATADIR%%/Gamez/cherrypy/test/test_xmlrpc.py
%%DATADIR%%/Gamez/cherrypy/test/webtest.py
%%DATADIR%%/Gamez/cherrypy/tutorial/README.txt
%%DATADIR%%/Gamez/cherrypy/tutorial/__init__.py
%%DATADIR%%/Gamez/cherrypy/tutorial/bonus-sqlobject.py
%%DATADIR%%/Gamez/cherrypy/tutorial/custom_error.html
%%DATADIR%%/Gamez/cherrypy/tutorial/pdf_file.pdf
%%DATADIR%%/Gamez/cherrypy/tutorial/tut01_helloworld.py
%%DATADIR%%/Gamez/cherrypy/tutorial/tut02_expose_methods.py
%%DATADIR%%/Gamez/cherrypy/tutorial/tut03_get_and_post.py
%%DATADIR%%/Gamez/cherrypy/tutorial/tut04_complex_site.py
%%DATADIR%%/Gamez/cherrypy/tutorial/tut05_derived_objects.py
%%DATADIR%%/Gamez/cherrypy/tutorial/tut06_default_method.py
%%DATADIR%%/Gamez/cherrypy/tutorial/tut07_sessions.py
%%DATADIR%%/Gamez/cherrypy/tutorial/tut08_generators_and_yield.py
%%DATADIR%%/Gamez/cherrypy/tutorial/tut09_files.py
%%DATADIR%%/Gamez/cherrypy/tutorial/tut10_http_errors.py
%%DATADIR%%/Gamez/cherrypy/tutorial/tutorial.conf
%%DATADIR%%/Gamez/cherrypy/wsgiserver/__init__.py
%%DATADIR%%/Gamez/cherrypy/wsgiserver/ssl_builtin.py
%%DATADIR%%/Gamez/cherrypy/wsgiserver/ssl_pyopenssl.py
%%DATADIR%%/Gamez/cherrypy/wsgiserver/wsgiserver2.py
%%DATADIR%%/Gamez/cherrypy/wsgiserver/wsgiserver3.py
%%DATADIR%%/Gamez/cron/__init__.py
%%DATADIR%%/Gamez/cron/base.py
%%DATADIR%%/Gamez/cron/eta.py
%%DATADIR%%/Gamez/cron/kinepolisrss.py
%%DATADIR%%/Gamez/cron/movierss.py
%%DATADIR%%/Gamez/cron/renamer.py
%%DATADIR%%/Gamez/cron/subtitle.py
%%DATADIR%%/Gamez/cron/trailer.py
%%DATADIR%%/Gamez/cron/trakt.py
%%DATADIR%%/Gamez/cron/yarr.py
%%DATADIR%%/Gamez/css/default/datatables.css
%%DATADIR%%/Gamez/css/default/datatables_images/Sorting icons.psd
%%DATADIR%%/Gamez/css/default/datatables_images/back_disabled.png
%%DATADIR%%/Gamez/css/default/datatables_images/back_enabled.png
%%DATADIR%%/Gamez/css/default/datatables_images/back_enabled_hover.png
%%DATADIR%%/Gamez/css/default/datatables_images/favicon.ico
%%DATADIR%%/Gamez/css/default/datatables_images/forward_disabled.png
%%DATADIR%%/Gamez/css/default/datatables_images/forward_enabled.png
%%DATADIR%%/Gamez/css/default/datatables_images/forward_enabled_hover.png
%%DATADIR%%/Gamez/css/default/datatables_images/shutdown.png
%%DATADIR%%/Gamez/css/default/datatables_images/sort_asc.png
%%DATADIR%%/Gamez/css/default/datatables_images/sort_asc_disabled.png
%%DATADIR%%/Gamez/css/default/datatables_images/sort_both.png
%%DATADIR%%/Gamez/css/default/datatables_images/sort_desc.png
%%DATADIR%%/Gamez/css/default/datatables_images/sort_desc_disabled.png
%%DATADIR%%/Gamez/css/default/jquery.ui.override.css
%%DATADIR%%/Gamez/css/default/navigation.css
%%DATADIR%%/Gamez/css/default/navigation_images/item-pointer-mover.gif
%%DATADIR%%/Gamez/css/default/navigation_images/item-pointer.gif
%%DATADIR%%/Gamez/css/default/navigation_images/lava.gif
%%DATADIR%%/Gamez/css/default/navigation_images/lava.png
%%DATADIR%%/Gamez/css/default/navigation_images/main-bg.png
%%DATADIR%%/Gamez/css/default/navigation_images/main-delimiter.png
%%DATADIR%%/Gamez/css/default/navigation_images/no-cover.jpg
%%DATADIR%%/Gamez/css/default/redmond/images/ui-bg_flat_0_aaaaaa_40x100.png
%%DATADIR%%/Gamez/css/default/redmond/images/ui-bg_flat_55_fbec88_40x100.png
%%DATADIR%%/Gamez/css/default/redmond/images/ui-bg_glass_75_d0e5f5_1x400.png
%%DATADIR%%/Gamez/css/default/redmond/images/ui-bg_glass_85_dfeffc_1x400.png
%%DATADIR%%/Gamez/css/default/redmond/images/ui-bg_glass_95_fef1ec_1x400.png
%%DATADIR%%/Gamez/css/default/redmond/images/ui-bg_gloss-wave_55_5c9ccc_500x100.png
%%DATADIR%%/Gamez/css/default/redmond/images/ui-bg_inset-hard_100_f5f8f9_1x100.png
%%DATADIR%%/Gamez/css/default/redmond/images/ui-bg_inset-hard_100_fcfdfd_1x100.png
%%DATADIR%%/Gamez/css/default/redmond/images/ui-icons_217bc0_256x240.png
%%DATADIR%%/Gamez/css/default/redmond/images/ui-icons_2e83ff_256x240.png
%%DATADIR%%/Gamez/css/default/redmond/images/ui-icons_469bdd_256x240.png
%%DATADIR%%/Gamez/css/default/redmond/images/ui-icons_6da8d5_256x240.png
%%DATADIR%%/Gamez/css/default/redmond/images/ui-icons_cd0a0a_256x240.png
%%DATADIR%%/Gamez/css/default/redmond/images/ui-icons_d8e7f3_256x240.png
%%DATADIR%%/Gamez/css/default/redmond/images/ui-icons_f9bd01_256x240.png
%%DATADIR%%/Gamez/css/default/redmond/jquery-ui-1.8.16.custom.css
%%DATADIR%%/Gamez/css/default/settings.css
%%DATADIR%%/Gamez/css/grey/datatables.css
%%DATADIR%%/Gamez/css/grey/datatables_images/Sorting icons.psd
%%DATADIR%%/Gamez/css/grey/datatables_images/back_disabled.png
%%DATADIR%%/Gamez/css/grey/datatables_images/back_enabled.png
%%DATADIR%%/Gamez/css/grey/datatables_images/back_enabled_hover.png
%%DATADIR%%/Gamez/css/grey/datatables_images/body-background.jpg
%%DATADIR%%/Gamez/css/grey/datatables_images/favicon.ico
%%DATADIR%%/Gamez/css/grey/datatables_images/forward_disabled.png
%%DATADIR%%/Gamez/css/grey/datatables_images/forward_enabled.png
%%DATADIR%%/Gamez/css/grey/datatables_images/forward_enabled_hover.png
%%DATADIR%%/Gamez/css/grey/datatables_images/shutdown.png
%%DATADIR%%/Gamez/css/grey/datatables_images/sort_asc.png
%%DATADIR%%/Gamez/css/grey/datatables_images/sort_asc_disabled.png
%%DATADIR%%/Gamez/css/grey/datatables_images/sort_both.png
%%DATADIR%%/Gamez/css/grey/datatables_images/sort_desc.png
%%DATADIR%%/Gamez/css/grey/datatables_images/sort_desc_disabled.png
%%DATADIR%%/Gamez/css/grey/jquery.ui.override.css
%%DATADIR%%/Gamez/css/grey/navigation.css
%%DATADIR%%/Gamez/css/grey/navigation_images/item-pointer-grey.gif
%%DATADIR%%/Gamez/css/grey/navigation_images/item-pointer-mover-grey.gif
%%DATADIR%%/Gamez/css/grey/navigation_images/lava-grey.gif
%%DATADIR%%/Gamez/css/grey/navigation_images/lava-grey.png
%%DATADIR%%/Gamez/css/grey/navigation_images/main-bg-grey.png
%%DATADIR%%/Gamez/css/grey/navigation_images/main-delimiter-grey.png
%%DATADIR%%/Gamez/css/grey/redmond/images/ui-bg_flat_0_aaaaaa_40x100.png
%%DATADIR%%/Gamez/css/grey/redmond/images/ui-bg_flat_55_fbec88_40x100.png
%%DATADIR%%/Gamez/css/grey/redmond/images/ui-bg_glass_75_d0e5f5_1x400.png
%%DATADIR%%/Gamez/css/grey/redmond/images/ui-bg_glass_85_dfeffc_1x400.png
%%DATADIR%%/Gamez/css/grey/redmond/images/ui-bg_glass_95_fef1ec_1x400.png
%%DATADIR%%/Gamez/css/grey/redmond/images/ui-bg_gloss-wave_55_5c9ccc_500x100.png
%%DATADIR%%/Gamez/css/grey/redmond/images/ui-bg_inset-hard_100_f5f8f9_1x100.png
%%DATADIR%%/Gamez/css/grey/redmond/images/ui-bg_inset-hard_100_fcfdfd_1x100.png
%%DATADIR%%/Gamez/css/grey/redmond/images/ui-icons_217bc0_256x240.png
%%DATADIR%%/Gamez/css/grey/redmond/images/ui-icons_2e83ff_256x240.png
%%DATADIR%%/Gamez/css/grey/redmond/images/ui-icons_469bdd_256x240.png
%%DATADIR%%/Gamez/css/grey/redmond/images/ui-icons_6da8d5_256x240.png
%%DATADIR%%/Gamez/css/grey/redmond/images/ui-icons_cd0a0a_256x240.png
%%DATADIR%%/Gamez/css/grey/redmond/images/ui-icons_d8e7f3_256x240.png
%%DATADIR%%/Gamez/css/grey/redmond/images/ui-icons_f9bd01_256x240.png
%%DATADIR%%/Gamez/css/grey/redmond/jquery-ui-1.8.16.custom.css
%%DATADIR%%/Gamez/css/grey/settings.css
%%DATADIR%%/Gamez/css/images/favicon.ico
%%DATADIR%%/Gamez/css/images/gamezlogo.png
%%DATADIR%%/Gamez/css/images/icon.delete.png
%%DATADIR%%/Gamez/css/images/icon.folder.gif
%%DATADIR%%/Gamez/css/images/icon.refresh.png
%%DATADIR%%/Gamez/css/images/icon.search.png
%%DATADIR%%/Gamez/css/images/logo.ico
%%DATADIR%%/Gamez/css/images/setting_downloader.jpg
%%DATADIR%%/Gamez/css/images/setting_folders.jpg
%%DATADIR%%/Gamez/css/images/setting_gamez.jpg
%%DATADIR%%/Gamez/css/images/setting_notifications.jpg
%%DATADIR%%/Gamez/css/images/setting_search.jpg
%%DATADIR%%/Gamez/css/js/jquery-1.6.2.min.js
%%DATADIR%%/Gamez/css/js/jquery-ui-1.8.16.custom.min.js
%%DATADIR%%/Gamez/css/js/jquery.dataTables.min.js
%%DATADIR%%/Gamez/css/js/menu.js
%%DATADIR%%/Gamez/gamez/ConfigFunctions.py
%%DATADIR%%/Gamez/gamez/Constants.py
%%DATADIR%%/Gamez/gamez/DBFunctions.py
%%DATADIR%%/Gamez/gamez/FolderFunctions.py
%%DATADIR%%/Gamez/gamez/GameTasks.py
%%DATADIR%%/Gamez/gamez/Helper.py
%%DATADIR%%/Gamez/gamez/Logger.py
%%DATADIR%%/Gamez/gamez/NFOCreator.py
%%DATADIR%%/Gamez/gamez/Notifications.py
%%DATADIR%%/Gamez/gamez/Postproccess.py
%%DATADIR%%/Gamez/gamez/TheGamesDBSearcher.py
%%DATADIR%%/Gamez/gamez/UpgradeFunctions.py
%%DATADIR%%/Gamez/gamez/WebRoot.py
%%DATADIR%%/Gamez/gamez/__init__.py
%%DATADIR%%/Gamez/gamez/xbmc/Xbmc.py
%%DATADIR%%/Gamez/gamez/xbmc/__init__.py
%%DATADIR%%/Gamez/lib/__init__.py
%%DATADIR%%/Gamez/lib/abgx360ap/BeautifulSoup.py
%%DATADIR%%/Gamez/lib/abgx360ap/README
%%DATADIR%%/Gamez/lib/abgx360ap/__init__.py
%%DATADIR%%/Gamez/lib/abgx360ap/abgx360.py
%%DATADIR%%/Gamez/lib/certgen.py
%%DATADIR%%/Gamez/lib/feedparser.py
%%DATADIR%%/Gamez/lib/gntp/__init__.py
%%DATADIR%%/Gamez/lib/gntp/notifier.py
%%DATADIR%%/Gamez/lib/notifo/Notifo.py
%%DATADIR%%/Gamez/lib/notifo/__init__.py
%%DATADIR%%/Gamez/lib/prowlpy/__init__.py
%%DATADIR%%/Gamez/lib/prowlpy/prowlpy.py
%%DATADIR%%/Gamez/lib/pynma/__init__.py
%%DATADIR%%/Gamez/lib/pynma/pynma.py
%%DATADIR%%/Gamez/postprocess/gamezPostProcess.py
%%DATADIR%%/Gamez/sample.init.ubuntu.txt
%%DATADIR%%/Gamez/tools/SETUP.md
%%DATADIR%%/Gamez/tools/setup.py
@unexec rm -r %D/share/gamez
@dirrmtry %%ETCDIR%%
@exec mkdir -p %D/etc/gamez
@exec chown -R media:media %D/etc/gamez  %D/share/gamez/*
