# ports/net-p2p/btsync/pkg-plist
@unexec %D/etc/rc.d/btsync forcestop 2>/dev/null || true
@unexec if cmp -s %D/etc/btsync/config.ini.sample %D/etc/btsync/config.ini; then rm -f %D/etc/btsync/config.ini; fi
%%ETCDIR%%/config.ini.sample
@dirrmtry %%ETCDIR%%
@exec if [ ! -f %D/etc/btsync/config.ini ] ; then cp -p %D/etc/btsync/config.ini.sample %D/etc/btsync/config.ini; fi
@exec chown -R btsync:btsync %D/etc/btsync
bin/btsync
