# ports/net-p2p/btsync/files/btsync.in
#!/bin/sh
#
# PROVIDE: %%PORTNAME%%
# REQUIRE: NETWORKING SERVERS DAEMON ldconfig resolv
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# btsync_enable (bool):	Set to NO by default.
#			Set it to YES to enable it.
# btsync_data_dir:	Directory where %%PORTNAME%% configuration
#			data is stored.
#			Default: %%PREFIX%%/etc/%%PORTNAME%%
# btsync_user:		The user account %%PORTNAME%% daemon runs as what
#			you want it to be. It uses '%%USER%%' user by
#			default. Do not sets it as empty or it will run
#			as root.
# btsync_group:		The group account %%PORTNAME%% daemon runs as what
#			you want it to be. It uses '%%GROUP%%' group by
#			default. Do not sets it as empty or it will run
#			as wheel.

. /etc/rc.subr
name="%%PORTNAME%%"
rcvar="${name}_enable"
load_rc_config $name

: ${btsync_enable:="NO"}
: ${btsync_user:="%%USER%%"}
: ${btsyncr_group:="%%GROUP%%"}
: ${btsync_data_dir:="%%PREFIX%%/etc/%%PORTNAME%%"}

pidfile="$btsync_data_dir/sync.pid"
command="%%PREFIX%%/bin/btsync"
command_args="--config $btsync_data_dir/config.ini"

run_rc_command "$1"
