# ports/games/gamez/files/gamez.in
#!/bin/sh
#
# PROVIDE: %%PORTNAME%%
# REQUIRE: NETWORKING SERVERS DAEMON ldconfig resolv
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# gamez_enable (bool):	Set to NO by default.
#			Set it to YES to enable it.
# gamez_data_dir:	Directory where %%PORTNAME%% configuration
#			data is stored.
#			Default: %%PREFIX%%/etc/%%PORTNAME%%
# gamez_user:		The user account %%PORTNAME%% daemon runs as what
#			you want it to be. It uses 'media' user by
#			default. Do not sets it as empty or it will run
#			as root.
# gamez_group:		The group account %%PORTNAME%% daemon runs as what
#			you want it to be. It uses 'media' group by
#			default. Do not sets it as empty or it will run
#			as wheel.

. /etc/rc.subr
name="%%PORTNAME%%"
rcvar="${name}_enable"
load_rc_config $name

: ${gamez_enable:="NO"}
: ${gamez_user:="%%USER%%"}
: ${gamez_group:="%%GROUP%%"}
: ${gamez_dir:="%%PREFIX%%/share/%%PORTNAME%%"}
: ${gamez_data_dir:="%%PREFIX%%/etc/%%PORTNAME%%"}

pidfile="$gamez_data_dir/%%PORTNAME%%.pid"
command="%%PYTHON_CMD%%"
command_args="$gamez_dir/Gamez/Gamez.py -c $gamez_data_dir/config.ini -b $gamez_data_dir -d -p $pidfile -n"

run_rc_command "$1"
