# ports/multimedia/maraschino/files/maraschino.in
#!/bin/sh
#
# PROVIDE: %%PORTNAME%%
# REQUIRE: NETWORKING SERVERS DAEMON ldconfig resolv
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# %%PORTNAME%%_enable (bool):	Set to NO by default.
#			Set it to YES to enable it.
# %%PORTNAME%%_data_dir:	Directory where %%PORTNAME%% configuration
#			data is stored.
#			Default: %%PREFIX%%/etc/%%PORTNAME%%
# %%PORTNAME%%_user:	The user account %%PORTNAME%% daemon runs as what
#			you want it to be. It uses '%%USER%%' user by
#			default. Do not sets it as empty or it will run
#			as root.
# %%PORTNAME%%_group:	The group account %%PORTNAME%% daemon runs as what
#			you want it to be. It uses '%%GROUP%%' group by
#			default. Do not sets it as empty or it will run
#			as wheel.

. /etc/rc.subr
name="%%PORTNAME%%"
rcvar="${name}_enable"
load_rc_config $name

: ${%%PORTNAME%%_enable:="NO"}
: ${%%PORTNAME%%_user:="%%USER%%"}
: ${%%PORTNAME%%_group:="%%GROUP%%"}
: ${%%PORTNAME%%_dir:="%%PREFIX%%/share/%%PORTNAME%%"}
: ${%%PORTNAME%%_data_dir:="%%PREFIX%%/etc/%%PORTNAME%%"}

pidfile="$maraschino_data_dir/%%PORTNAME%%.pid"
command="%%PYTHON_CMD%%"
command_args="$maraschino_dir/maraschino/Maraschino.py --datadir=$maraschino_data_dir -d --pidfile=$pidfile"

run_rc_command "$1"
