# ports/deskutils/lazylibrarian/pkg-plist
@unexec %D/etc/rc.d/lazylibrarian forcestop 2>/dev/null || true
%%DATADIR%%/LazyLibrarian/LazyLibrarian.app/Contents/Info.plist
%%DATADIR%%/LazyLibrarian/LazyLibrarian.app/Contents/MacOS/applet
%%DATADIR%%/LazyLibrarian/LazyLibrarian.app/Contents/PkgInfo
%%DATADIR%%/LazyLibrarian/LazyLibrarian.app/Contents/Resources/Scripts/main.scpt
%%DATADIR%%/LazyLibrarian/LazyLibrarian.app/Contents/Resources/applet.icns
%%DATADIR%%/LazyLibrarian/LazyLibrarian.app/Contents/Resources/applet.rsrc
%%DATADIR%%/LazyLibrarian/LazyLibrarian.app/Contents/Resources/description.rtfd/TXT.rtf
%%DATADIR%%/LazyLibrarian/LazyLibrarian.py
%%DATADIR%%/LazyLibrarian/README.md
%%DATADIR%%/LazyLibrarian/UNITTESTING.md
%%DATADIR%%/LazyLibrarian/branch-summary.log
%%DATADIR%%/LazyLibrarian/cherrypy/LICENSE.txt
%%DATADIR%%/LazyLibrarian/cherrypy/__init__.py
%%DATADIR%%/LazyLibrarian/cherrypy/_cpchecker.py
%%DATADIR%%/LazyLibrarian/cherrypy/_cpcompat.py
%%DATADIR%%/LazyLibrarian/cherrypy/_cpconfig.py
%%DATADIR%%/LazyLibrarian/cherrypy/_cpdispatch.py
%%DATADIR%%/LazyLibrarian/cherrypy/_cperror.py
%%DATADIR%%/LazyLibrarian/cherrypy/_cplogging.py
%%DATADIR%%/LazyLibrarian/cherrypy/_cpmodpy.py
%%DATADIR%%/LazyLibrarian/cherrypy/_cpnative_server.py
%%DATADIR%%/LazyLibrarian/cherrypy/_cpreqbody.py
%%DATADIR%%/LazyLibrarian/cherrypy/_cprequest.py
%%DATADIR%%/LazyLibrarian/cherrypy/_cpserver.py
%%DATADIR%%/LazyLibrarian/cherrypy/_cpthreadinglocal.py
%%DATADIR%%/LazyLibrarian/cherrypy/_cptools.py
%%DATADIR%%/LazyLibrarian/cherrypy/_cptree.py
%%DATADIR%%/LazyLibrarian/cherrypy/_cpwsgi.py
%%DATADIR%%/LazyLibrarian/cherrypy/_cpwsgi_server.py
%%DATADIR%%/LazyLibrarian/cherrypy/cherryd
%%DATADIR%%/LazyLibrarian/cherrypy/favicon.ico
%%DATADIR%%/LazyLibrarian/cherrypy/lib/__init__.py
%%DATADIR%%/LazyLibrarian/cherrypy/lib/auth.py
%%DATADIR%%/LazyLibrarian/cherrypy/lib/auth_basic.py
%%DATADIR%%/LazyLibrarian/cherrypy/lib/auth_digest.py
%%DATADIR%%/LazyLibrarian/cherrypy/lib/caching.py
%%DATADIR%%/LazyLibrarian/cherrypy/lib/covercp.py
%%DATADIR%%/LazyLibrarian/cherrypy/lib/cpstats.py
%%DATADIR%%/LazyLibrarian/cherrypy/lib/cptools.py
%%DATADIR%%/LazyLibrarian/cherrypy/lib/encoding.py
%%DATADIR%%/LazyLibrarian/cherrypy/lib/http.py
%%DATADIR%%/LazyLibrarian/cherrypy/lib/httpauth.py
%%DATADIR%%/LazyLibrarian/cherrypy/lib/httputil.py
%%DATADIR%%/LazyLibrarian/cherrypy/lib/jsontools.py
%%DATADIR%%/LazyLibrarian/cherrypy/lib/profiler.py
%%DATADIR%%/LazyLibrarian/cherrypy/lib/reprconf.py
%%DATADIR%%/LazyLibrarian/cherrypy/lib/sessions.py
%%DATADIR%%/LazyLibrarian/cherrypy/lib/static.py
%%DATADIR%%/LazyLibrarian/cherrypy/lib/xmlrpc.py
%%DATADIR%%/LazyLibrarian/cherrypy/process/__init__.py
%%DATADIR%%/LazyLibrarian/cherrypy/process/plugins.py
%%DATADIR%%/LazyLibrarian/cherrypy/process/servers.py
%%DATADIR%%/LazyLibrarian/cherrypy/process/win32.py
%%DATADIR%%/LazyLibrarian/cherrypy/process/wspbus.py
%%DATADIR%%/LazyLibrarian/cherrypy/scaffold/__init__.py
%%DATADIR%%/LazyLibrarian/cherrypy/scaffold/apache-fcgi.conf
%%DATADIR%%/LazyLibrarian/cherrypy/scaffold/example.conf
%%DATADIR%%/LazyLibrarian/cherrypy/scaffold/site.conf
%%DATADIR%%/LazyLibrarian/cherrypy/scaffold/static/made_with_cherrypy_small.png
%%DATADIR%%/LazyLibrarian/cherrypy/test/__init__.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/_test_decorators.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/_test_states_demo.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/benchmark.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/checkerdemo.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/fastcgi.conf
%%DATADIR%%/LazyLibrarian/cherrypy/test/fcgi.conf
%%DATADIR%%/LazyLibrarian/cherrypy/test/helper.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/logtest.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/modfastcgi.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/modfcgid.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/modpy.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/modwsgi.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/native-server.ini
%%DATADIR%%/LazyLibrarian/cherrypy/test/sessiondemo.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/static/dirback.jpg
%%DATADIR%%/LazyLibrarian/cherrypy/test/static/index.html
%%DATADIR%%/LazyLibrarian/cherrypy/test/style.css
%%DATADIR%%/LazyLibrarian/cherrypy/test/test.pem
%%DATADIR%%/LazyLibrarian/cherrypy/test/test_auth_basic.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/test_auth_digest.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/test_bus.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/test_caching.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/test_config.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/test_config_server.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/test_conn.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/test_core.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/test_dynamicobjectmapping.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/test_encoding.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/test_etags.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/test_http.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/test_httpauth.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/test_httplib.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/test_json.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/test_logging.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/test_mime.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/test_misc_tools.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/test_objectmapping.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/test_proxy.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/test_refleaks.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/test_request_obj.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/test_routes.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/test_session.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/test_sessionauthenticate.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/test_states.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/test_static.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/test_tools.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/test_tutorials.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/test_virtualhost.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/test_wsgi_ns.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/test_wsgi_vhost.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/test_wsgiapps.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/test_xmlrpc.py
%%DATADIR%%/LazyLibrarian/cherrypy/test/webtest.py
%%DATADIR%%/LazyLibrarian/cherrypy/wsgiserver/__init__.py
%%DATADIR%%/LazyLibrarian/cherrypy/wsgiserver/ssl_builtin.py
%%DATADIR%%/LazyLibrarian/cherrypy/wsgiserver/ssl_pyopenssl.py
%%DATADIR%%/LazyLibrarian/data/css/data_table.css
%%DATADIR%%/LazyLibrarian/data/css/ipad_landscape.css
%%DATADIR%%/LazyLibrarian/data/css/ipad_portrait.css
%%DATADIR%%/LazyLibrarian/data/css/iphone4_portrait.css
%%DATADIR%%/LazyLibrarian/data/css/mobile.css
%%DATADIR%%/LazyLibrarian/data/css/style.css
%%DATADIR%%/LazyLibrarian/data/images/0-stars.png
%%DATADIR%%/LazyLibrarian/data/images/1-stars.png
%%DATADIR%%/LazyLibrarian/data/images/2-stars.png
%%DATADIR%%/LazyLibrarian/data/images/3-stars.png
%%DATADIR%%/LazyLibrarian/data/images/4-stars.png
%%DATADIR%%/LazyLibrarian/data/images/5-stars.png
%%DATADIR%%/LazyLibrarian/data/images/back_disabled.jpg
%%DATADIR%%/LazyLibrarian/data/images/back_enabled.jpg
%%DATADIR%%/LazyLibrarian/data/images/checkmark.png
%%DATADIR%%/LazyLibrarian/data/images/favicon.ico
%%DATADIR%%/LazyLibrarian/data/images/forward_disabled.jpg
%%DATADIR%%/LazyLibrarian/data/images/forward_enabled.jpg
%%DATADIR%%/LazyLibrarian/data/images/headphoneslogo.png
%%DATADIR%%/LazyLibrarian/data/images/icon_search.png
%%DATADIR%%/LazyLibrarian/data/images/lazylibrarian.png
%%DATADIR%%/LazyLibrarian/data/images/loader_black.gif
%%DATADIR%%/LazyLibrarian/data/images/nav/authors.jpg
%%DATADIR%%/LazyLibrarian/data/images/nav/authors.png
%%DATADIR%%/LazyLibrarian/data/images/nav/authorsActive.jpg
%%DATADIR%%/LazyLibrarian/data/images/nav/authorsActive.png
%%DATADIR%%/LazyLibrarian/data/images/nav/books.jpg
%%DATADIR%%/LazyLibrarian/data/images/nav/books.png
%%DATADIR%%/LazyLibrarian/data/images/nav/booksActive.jpg
%%DATADIR%%/LazyLibrarian/data/images/nav/booksActive.png
%%DATADIR%%/LazyLibrarian/data/images/nav/config.jpg
%%DATADIR%%/LazyLibrarian/data/images/nav/config.png
%%DATADIR%%/LazyLibrarian/data/images/nav/configActive.jpg
%%DATADIR%%/LazyLibrarian/data/images/nav/configActive.png
%%DATADIR%%/LazyLibrarian/data/images/nav/logs.jpg
%%DATADIR%%/LazyLibrarian/data/images/nav/logs.png
%%DATADIR%%/LazyLibrarian/data/images/nav/logsActive.jpg
%%DATADIR%%/LazyLibrarian/data/images/nav/logsActive.png
%%DATADIR%%/LazyLibrarian/data/images/nocover.png
%%DATADIR%%/LazyLibrarian/data/images/nostar.png
%%DATADIR%%/LazyLibrarian/data/images/restart.png
%%DATADIR%%/LazyLibrarian/data/images/shutdown.png
%%DATADIR%%/LazyLibrarian/data/images/sprite.png
%%DATADIR%%/LazyLibrarian/data/images/sprite1.png
%%DATADIR%%/LazyLibrarian/data/images/star.png
%%DATADIR%%/LazyLibrarian/data/images/templogo.png
%%DATADIR%%/LazyLibrarian/data/interfaces/default/album.html
%%DATADIR%%/LazyLibrarian/data/interfaces/default/author.html
%%DATADIR%%/LazyLibrarian/data/interfaces/default/base.html
%%DATADIR%%/LazyLibrarian/data/interfaces/default/books.html
%%DATADIR%%/LazyLibrarian/data/interfaces/default/config.html
%%DATADIR%%/LazyLibrarian/data/interfaces/default/extras.html
%%DATADIR%%/LazyLibrarian/data/interfaces/default/history.html
%%DATADIR%%/LazyLibrarian/data/interfaces/default/index.html
%%DATADIR%%/LazyLibrarian/data/interfaces/default/logs.html
%%DATADIR%%/LazyLibrarian/data/interfaces/default/manageartists.html
%%DATADIR%%/LazyLibrarian/data/interfaces/default/managenew.html
%%DATADIR%%/LazyLibrarian/data/interfaces/default/searchresults.html
%%DATADIR%%/LazyLibrarian/data/interfaces/default/shutdown.html
%%DATADIR%%/LazyLibrarian/data/interfaces/default/upcoming.html
%%DATADIR%%/LazyLibrarian/data/js/adapt.js
%%DATADIR%%/LazyLibrarian/data/js/adapt.min.js
%%DATADIR%%/LazyLibrarian/data/js/jquery-1.8.2.js
%%DATADIR%%/LazyLibrarian/data/js/jquery.clearsearch.js
%%DATADIR%%/LazyLibrarian/data/js/lazylibrarian.js
%%DATADIR%%/LazyLibrarian/data/js/libs/dd_belatedpng.js
%%DATADIR%%/LazyLibrarian/data/js/libs/jquery-1.6.2.js
%%DATADIR%%/LazyLibrarian/data/js/libs/jquery-1.6.2.min.js
%%DATADIR%%/LazyLibrarian/data/js/libs/jquery.dataTables.min.js
%%DATADIR%%/LazyLibrarian/data/js/libs/modernizr-1.7.min.js
%%DATADIR%%/LazyLibrarian/data/js/mootools.js
%%DATADIR%%/LazyLibrarian/data/js/mootools_more.js
%%DATADIR%%/LazyLibrarian/data/js/mootools_tween_css3.js
%%DATADIR%%/LazyLibrarian/data/js/nav.js
%%DATADIR%%/LazyLibrarian/data/js/plugins.js
%%DATADIR%%/LazyLibrarian/data/js/question.js
%%DATADIR%%/LazyLibrarian/data/js/script.js
%%DATADIR%%/LazyLibrarian/init/COPYING.txt
%%DATADIR%%/LazyLibrarian/init/freebsd.initd
%%DATADIR%%/LazyLibrarian/init/ubuntu.default
%%DATADIR%%/LazyLibrarian/init/ubuntu.initd
%%DATADIR%%/LazyLibrarian/lazylibrarian/SimpleCache.py
%%DATADIR%%/LazyLibrarian/lazylibrarian/__init__.py
%%DATADIR%%/LazyLibrarian/lazylibrarian/database.py
%%DATADIR%%/LazyLibrarian/lazylibrarian/formatter.py
%%DATADIR%%/LazyLibrarian/lazylibrarian/gb.py
%%DATADIR%%/LazyLibrarian/lazylibrarian/gr.py
%%DATADIR%%/LazyLibrarian/lazylibrarian/importer.py
%%DATADIR%%/LazyLibrarian/lazylibrarian/logger.py
%%DATADIR%%/LazyLibrarian/lazylibrarian/postprocess.py
%%DATADIR%%/LazyLibrarian/lazylibrarian/providers.py
%%DATADIR%%/LazyLibrarian/lazylibrarian/sabnzbd.py
%%DATADIR%%/LazyLibrarian/lazylibrarian/searchnzb.py
%%DATADIR%%/LazyLibrarian/lazylibrarian/unittests/test_searchnzb.py
%%DATADIR%%/LazyLibrarian/lazylibrarian/unittests/test_versioncheck.py
%%DATADIR%%/LazyLibrarian/lazylibrarian/version.py
%%DATADIR%%/LazyLibrarian/lazylibrarian/versioncheck.py
%%DATADIR%%/LazyLibrarian/lazylibrarian/webServe.py
%%DATADIR%%/LazyLibrarian/lazylibrarian/webStart.py
%%DATADIR%%/LazyLibrarian/lib/__init__.py
%%DATADIR%%/LazyLibrarian/lib/apscheduler/__init__.py
%%DATADIR%%/LazyLibrarian/lib/apscheduler/events.py
%%DATADIR%%/LazyLibrarian/lib/apscheduler/job.py
%%DATADIR%%/LazyLibrarian/lib/apscheduler/jobstores/__init__.py
%%DATADIR%%/LazyLibrarian/lib/apscheduler/jobstores/base.py
%%DATADIR%%/LazyLibrarian/lib/apscheduler/jobstores/mongodb_store.py
%%DATADIR%%/LazyLibrarian/lib/apscheduler/jobstores/ram_store.py
%%DATADIR%%/LazyLibrarian/lib/apscheduler/jobstores/shelve_store.py
%%DATADIR%%/LazyLibrarian/lib/apscheduler/jobstores/sqlalchemy_store.py
%%DATADIR%%/LazyLibrarian/lib/apscheduler/scheduler.py
%%DATADIR%%/LazyLibrarian/lib/apscheduler/threadpool.py
%%DATADIR%%/LazyLibrarian/lib/apscheduler/triggers/__init__.py
%%DATADIR%%/LazyLibrarian/lib/apscheduler/triggers/cron/__init__.py
%%DATADIR%%/LazyLibrarian/lib/apscheduler/triggers/cron/expressions.py
%%DATADIR%%/LazyLibrarian/lib/apscheduler/triggers/cron/fields.py
%%DATADIR%%/LazyLibrarian/lib/apscheduler/triggers/interval.py
%%DATADIR%%/LazyLibrarian/lib/apscheduler/triggers/simple.py
%%DATADIR%%/LazyLibrarian/lib/apscheduler/util.py
%%DATADIR%%/LazyLibrarian/lib/configobj.py
%%DATADIR%%/LazyLibrarian/lib/dogpile/__init__.py
%%DATADIR%%/LazyLibrarian/lib/dogpile/cache/__init__.py
%%DATADIR%%/LazyLibrarian/lib/dogpile/cache/api.py
%%DATADIR%%/LazyLibrarian/lib/dogpile/cache/backends/__init__.py
%%DATADIR%%/LazyLibrarian/lib/dogpile/cache/backends/file.py
%%DATADIR%%/LazyLibrarian/lib/dogpile/cache/backends/memcached.py
%%DATADIR%%/LazyLibrarian/lib/dogpile/cache/backends/memory.py
%%DATADIR%%/LazyLibrarian/lib/dogpile/cache/backends/redis.py
%%DATADIR%%/LazyLibrarian/lib/dogpile/cache/compat.py
%%DATADIR%%/LazyLibrarian/lib/dogpile/cache/plugins/__init__.py
%%DATADIR%%/LazyLibrarian/lib/dogpile/cache/plugins/mako_cache.py
%%DATADIR%%/LazyLibrarian/lib/dogpile/cache/region.py
%%DATADIR%%/LazyLibrarian/lib/dogpile/cache/util.py
%%DATADIR%%/LazyLibrarian/lib/fuzzywuzzy/StringMatcher.py
%%DATADIR%%/LazyLibrarian/lib/fuzzywuzzy/__init__.py
%%DATADIR%%/LazyLibrarian/lib/fuzzywuzzy/fuzz.py
%%DATADIR%%/LazyLibrarian/lib/fuzzywuzzy/process.py
%%DATADIR%%/LazyLibrarian/lib/fuzzywuzzy/utils.py
%%DATADIR%%/LazyLibrarian/lib/simplejson/__init__.py
%%DATADIR%%/LazyLibrarian/lib/simplejson/_speedups.c
%%DATADIR%%/LazyLibrarian/lib/simplejson/decoder.py
%%DATADIR%%/LazyLibrarian/lib/simplejson/encoder.py
%%DATADIR%%/LazyLibrarian/lib/simplejson/scanner.py
%%DATADIR%%/LazyLibrarian/lib/validate.py
%%DATADIR%%/LazyLibrarian/mako/LICENSE
%%DATADIR%%/LazyLibrarian/mako/__init__.py
%%DATADIR%%/LazyLibrarian/mako/_ast_util.py
%%DATADIR%%/LazyLibrarian/mako/ast.py
%%DATADIR%%/LazyLibrarian/mako/cache.py
%%DATADIR%%/LazyLibrarian/mako/codegen.py
%%DATADIR%%/LazyLibrarian/mako/exceptions.py
%%DATADIR%%/LazyLibrarian/mako/ext/__init__.py
%%DATADIR%%/LazyLibrarian/mako/ext/autohandler.py
%%DATADIR%%/LazyLibrarian/mako/ext/babelplugin.py
%%DATADIR%%/LazyLibrarian/mako/ext/preprocessors.py
%%DATADIR%%/LazyLibrarian/mako/ext/pygmentplugin.py
%%DATADIR%%/LazyLibrarian/mako/ext/turbogears.py
%%DATADIR%%/LazyLibrarian/mako/filters.py
%%DATADIR%%/LazyLibrarian/mako/lexer.py
%%DATADIR%%/LazyLibrarian/mako/lookup.py
%%DATADIR%%/LazyLibrarian/mako/parsetree.py
%%DATADIR%%/LazyLibrarian/mako/pygen.py
%%DATADIR%%/LazyLibrarian/mako/pyparser.py
%%DATADIR%%/LazyLibrarian/mako/runtime.py
%%DATADIR%%/LazyLibrarian/mako/template.py
%%DATADIR%%/LazyLibrarian/mako/util.py
@unexec rm -r %D/share/lazylibrarian
@dirrmtry %%ETCDIR%%
@exec mkdir -p %D/etc/lazylibrarian
@exec chown -R media:media %D/etc/lazylibrarian %D/share/lazylibrarian/*
