# ports/multimedia/htpc-manager/files/htpc-manager.in
#!/bin/sh
#
# PROVIDE: %%PORTNAME%%
# REQUIRE: NETWORKING SERVERS DAEMON ldconfig resolv
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# htpc_manager_enable (bool):	Set to NO by default.
#			Set it to YES to enable it.
# htpc_manager_user:	The user account %%PORTNAME%% daemon runs as what
#			you want it to be. It uses '%%USER%%' user by
#			default. Do not sets it as empty or it will run
#			as root.
# htpc_manager_group:	The group account %%PORTNAME%% daemon runs as what
#			you want it to be. It uses '%%GROUP%%' group by
#			default. Do not sets it as empty or it will run
#			as wheel.
# htpc_manager_data_dir:	Directory where %%PORTNAME%% configuration
#				data is stored.
#				Default: /var/db/%%PORTNAME%%

. /etc/rc.subr
name="htpc_manager"
rcvar="${name}_enable"
load_rc_config $name

: ${htpc_manager_enable:="NO"}
: ${htpc_manager_user:="%%USER%%"}
: ${htpc_manager_group:="%%GROUP%%"}
: ${htpc_manager_data_dir:="/var/db/%%PORTNAME%%"}

pidfile="$htpc_manager_data_dir/%%PORTNAME%%.pid"
command="%%PYTHON_CMD%%"
command_args="%%DATADIR%%/HTPC-Manager/Htpc.py --datadir $htpc_manager_data_dir --daemon --pid $pidfile"

run_rc_command "$1"
