# ports/news/sabnzbd/files/sabnzbd.in
#!/bin/sh
#
# PROVIDE: %%PORTNAME%%
# REQUIRE: NETWORKING SERVERS DAEMON ldconfig resolv
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# %%PORTNAME%%_enable (bool):	Set to NO by default.
# 				Set it to YES to enable it.
# %%PORTNAME%%_data_dir:	Directory where %%PORTNAME%% configuration
# 				data is stored.
# 				Default: %%PREFIX%%/etc/%%PORTNAME%%
# %%PORTNAME%%_user:		The user account %%PORTNAME%% daemon runs as what
# 				you want it to be. It uses '%%USER%%' user by
# 				default. Do not sets it as empty or it will run
# 				as root.
# %%PORTNAME%%_group:		The group account %%PORTNAME%% daemon runs as what
# 				you want it to be. It uses '%%GROUP%%' group by
# 				default. Do not sets it as empty or it will run
# 				as wheel.

. /etc/rc.subr
name="%%PORTNAME%%"
rcvar="${name}_enable"
load_rc_config $name

: ${%%PORTNAME%%_enable:="NO"}
: ${%%PORTNAME%%_user:="%%USER%%"}
: ${%%PORTNAME%%_group:="%%GROUP%%"}
: ${%%PORTNAME%%_dir:="%%PREFIX%%/share/%%PORTNAME%%"}
: ${%%PORTNAME%%_data_dir:="%%PREFIX%%/etc/%%PORTNAME%%"}

pidfile="$sabnzbd_data_dir/%%PORTNAME%%.pid"
command="%%PYTHON_CMD%%"
command_args="$sabnzbd_dir/SABnzbd/SABnzbd.py -s 0.0.0.0 -f $sabnzbd_data_dir/${name}.ini --pidfile $pidfile -d -b 0"
extra_commands="update"
update_cmd="%%PORTNAME%%_update"

%%PORTNAME%%_update() {
	if [ -n "$rc_pid" ]; then
		run_rc_command onestop
		running=true
	else
		running=false
	fi
	fetch -o $sabnzbd_dir/download.tar.gz http://sourceforge.net/projects/sabnzbdplus/files/latest/download
	tar -xvf $sabnzbd_dir/download.tar.gz -C $sabnzbd_dir
	if [ $? -eq 0 ]; then		
		rm -r $sabnzbd_dir/download.tar.gz $sabnzbd_dir/SABnzbd
		mv $sabnzbd_dir/SABnzbd-* $sabnzbd_dir/SABnzbd
		chown -R $sabnzbd_user:$sabnzbd_group $sabnzbd_dir/SABnzbd
	fi
	if $running; then
		run_rc_command onestart
	fi
}

run_rc_command "$1"
