# ports/multimedia/xdm/files/xdm.in
#!/bin/sh
#
# PROVIDE: xdm
# REQUIRE: NETWORKING SERVERS DAEMON ldconfig resolv
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# xdm_enable (bool):	Set to NO by default.
#			Set it to YES to enable it.
# xdm_data_dir:		Directory where xdm configuration
#			data is stored.
#			Default: %%PREFIX%%/etc/xdm
# xdm_user:		The user account xdm daemon runs as what
#			you want it to be. It uses 'media' user by
#			default. Do not sets it as empty or it will run
#			as root.
# xdm_group:		The group account xdm daemon runs as what
#			you want it to be. It uses 'media' group by
#			default. Do not sets it as empty or it will run
#			as wheel.

. /etc/rc.subr
name="xdm"
rcvar="${name}_enable"
load_rc_config $name

: ${xdm_enable:="NO"}
: ${xdm_user:="media"}
: ${xdm_group:="media"}
: ${xdm_dir:="%%PREFIX%%/share/$name"}
: ${xdm_data_dir:="%%PREFIX%%/etc/$name"}

pidfile="$xdm_data_dir/$name.pid"
command="%%PYTHON_CMD%%"
command_args="$xdm_dir/XDM/XDM.py -b $xdm_data_dir -p $pidfile -d -n"

run_rc_command "$1"
